# tags/afrimesh-1.0pre-r1153-unstable/package-scripts/openwrt/afrimesh/files/etc/crontabs/afrimesh
*/1 * * * * /www/cgi-bin/status 2> /dev/null 
